# display/gpdtouch-wake.service
[Unit]
Description=GPD Touchscreen Rotation Service on Wake
After=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gpdtouch service &

[Install]
WantedBy=multi-user.target
